# tender-watch/ops/com.openclaw.tender-watch.incremental.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openclaw.tender-watch.incremental</string>

  <key>WorkingDirectory</key>
  <string>/Users/openclaw/Documents/Playground/tender-watch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>./run_monitor.sh incremental &gt;&gt; logs/launchd_incremental.log 2&gt;&amp;1</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>15</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>18</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/openclaw/Documents/Playground/tender-watch/logs/launchd_stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/openclaw/Documents/Playground/tender-watch/logs/launchd_stderr.log</string>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
